# tempo-agent.service
[Unit]
Description=Tempo Agent Slack Bot
After=network.target docker.service

[Service]
User=ubuntu
Group=ubuntu
Type=simple
WorkingDirectory=/home/ubuntu/github/tempoxyz/ai/cli/agent
Environment=PATH=/home/ubuntu/.local/bin:/usr/bin
EnvironmentFile=/home/ubuntu/github/tempoxyz/ai/cli/agent/.env
ExecStart=/home/ubuntu/.local/bin/uv run tempo-agent run
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
